# competition/servers/pypy-cyclone-N.ini
[program:pypy-cyclone-N]
command=./pypy-2.0.2/bin/pypy wsdemo-cyclone-N.py
autostart=false
autorestart=false
startretries=0
stopasgroup=true
killasgroup=true
