# competition/servers/pypy-cyclone-1.ini
[program:pypy-cyclone-1]
command=./pypy-2.0.2/bin/pypy wsdemo-cyclone-1.py
autostart=false
autorestart=false
startretries=0
stopasgroup=true
killasgroup=true
