# competition/servers/python-cyclone-1.ini
[program:python-cyclone-1]
command=python wsdemo-cyclone-1.py
autostart=false
autorestart=false
startretries=0
stopasgroup=true
killasgroup=true
