# competition/servers/python-cyclone-N.ini
[program:python-cyclone-N]
command=python wsdemo-cyclone-N.py
autostart=false
autorestart=false
startretries=0
stopasgroup=true
killasgroup=true
