# competition/servers/vertx-java.ini
[program:vertx-java]
command=vertx run WsdemoVertx.java
autostart=false
autorestart=false
startretries=0
stopasgroup=true
killasgroup=true