# competition/servers/vertx-python.ini
[program:vertx-python]
command=vertx run WsdemoVertx.py
autostart=false
autorestart=false
startretries=0
stopasgroup=true
killasgroup=true